# infra/deployment/systemd/chad-core-api.service
[Unit]
Description=Chad-Core API Service
After=network.target

[Service]
Type=simple
User=chad-core
WorkingDirectory=/opt/chad-core
Environment="PATH=/opt/chad-core/venv/bin"
ExecStart=/opt/chad-core/venv/bin/uvicorn apps.core_api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
